# winwinlinux/config/includes.chroot/etc/cron.hourly/winwin-tools
#!/bin/sh

/usr/bin/ww-cron-root 2>/dev/null
